# launchd/com.mitchelldynamics.caddy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mitchelldynamics.caddy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/caddy</string>
    <string>run</string>
    <string>--config</string>
    <string>/Users/Shared/server/config/caddy/Caddyfile</string>
    <string>--envfile</string>
    <string>/Users/Shared/server/config/caddy/.env</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/Shared/server/config/caddy</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/Shared/server/config/caddy</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/Shared/server/config/caddy/logs/caddy.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/Shared/server/config/caddy/logs/caddy.err.log</string>
</dict>
</plist>
